# monitor-restart-vpn/wg-monitor.service
[Unit]
Description=Monitor and Auto-Restart WireGuard
After=network.target

[Service]
ExecStart=/opt/wg-monitor.sh
