# monitor-restart-vpn/wg-monitor.timer
[Unit]
Description=Run WireGuard monitor every minute

[Timer]
OnUnitActiveSec=60
AccuracySec=1s

[Install]
WantedBy=timers.target
